# bin/assets/empoleos/empoleos.service
[Unit]
Description=Empoleos Auto Updates

[Service]
ExecStart=/bin/bash /etc/empoleos/init.sh

[Install]
WantedBy=network.target
