# bin/empoleos-init/empoleos-init.service
[Unit]
Description=Empoleos Init

[Service]
ExecStart=/bin/bash /etc/empoleos-init/run.sh
Nice=-15

[Install]
WantedBy=network.target
