# bin/empoleos-init/empoleos-init.service
[Unit]
Description=Empoleos Init

[Service]
ExecStart=/bin/bash /etc/empoleos-init/run.sh
StandardOutput=/home/$USER/empoleos-init.log
StandardError=/home/$USER/empoleos-init.log

[Install]
WantedBy=network.target
