# bin/empoleos-installer/empoleos-installer.service
[Unit]
Description=Empoleos Installer

[Service]
ExecStart=/bin/bash /etc/empoleos-installer/run.sh
Nice=-15

[Install]
WantedBy=network.target
